# docker/hack-chat/root/etc/services.d/client/run
#!/usr/bin/with-contenv bash

cd /app/client

exec \
	s6-setuidgid abc http-server