# docker/hack-chat/root/etc/services.d/server/run
#!/usr/bin/with-contenv bash

cd /app/server

exec \
	s6-setuidgid abc npm start --production